# SFC Orchestrator/VNFs/custom/custom.service
[Unit]
Description=User-defined custom VNF

[Service]
Type=simple
ExecStart=/opt/custom.sh

[Install]
WantedBy=multi-user.target
