# SFC Orchestrator/VNFs/firewall/firewall.service
[Unit]
Description=User-defined firewall VNF

[Service]
Type=simple
ExecStart=/opt/firewall.sh

[Install]
WantedBy=multi-user.target

