# SFC Orchestrator/VNFs/spamfilter/spamfilter.service
[Unit]
Description=User-defined Spamfilter VNF

[Service]
Type=simple
ExecStart=/opt/spamfilter.sh

[Install]
WantedBy=multi-user.target
