# SFC Orchestrator/VNFs/packetsniffer/packetsniffer.service
[Unit]
Description= User-defined PacketSniffer VNF

[Service]
Type=simple
ExecStart=/opt/packetsniffer.sh

[Install]
WantedBy=multi-user.target
