# SFC Orchestrator/VNFs/server/server.service
[Unit]
Description=User-defined Server VNF

[Service]
Type=simple
ExecStart=/opt/server.sh

[Install]
WantedBy=multi-user.target
