# SFC Orchestrator/VNFs/test/test.service
[Unit]
Description=User-defined Test VNF

[Service]
Type=simple
ExecStart=/opt/test.sh

[Install]
WantedBy=multi-user.target
