# SFC Orchestrator/VNFs/nat/nat.service
[Unit]
Description=User-defined NAT VNF

[Service]
Type=simple
ExecStart=/opt/nat.sh

[Install]
WantedBy=multi-user.target
